# EXCELBOT/telegram-bot.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/EXCELBOT
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/EXCELBOT/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

